# roles/web/deploy/files/daoistic.socket
[Unit]
Description=daoistic socket

[Socket]
ListenStream=/run/daoistic/socket

[Install]
WantedBy=sockets.target
