# role/deploy/files/daoistic-web.socket
[Unit]
Description=daoistic-web socket

[Socket]
ListenStream=/run/daoistic-web/socket

[Install]
WantedBy=sockets.target
